# supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── FastAPI backend ──────────────────────────────────────────────────────────
[program:fastapi]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app/backend
environment=PYTHONPATH="/app/backend",UPLOAD_DIR="/app/uploads"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi-err.log
stdout_logfile_maxbytes=10MB
priority=10

; ── Next.js frontend ─────────────────────────────────────────────────────────
[program:nextjs]
command=node_modules/.bin/next start --port 3000
directory=/app/frontend
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nextjs.log
stderr_logfile=/var/log/supervisor/nextjs-err.log
stdout_logfile_maxbytes=10MB
priority=20

; ── nginx reverse proxy ───────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-err.log
priority=30
